# deployment/supervisor.conf
# Soniq Supervisor Configuration
# Place this file in /etc/supervisor/conf.d/soniq.conf

\[group:soniq]
programs=soniq_worker,soniq_dashboard,soniq_scheduler

[program:soniq_worker]
command=/opt/soniq/venv/bin/soniq worker --concurrency=4
directory=/opt/soniq
user=soniq
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/soniq/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
# soniq_scheduler is configured separately below; suppress the
# @periodic-without-scheduler WARN to keep worker logs clean.
environment=SONIQ_DATABASE_URL="postgresql://soniq:password@localhost/soniq_prod",SONIQ_LOG_LEVEL="INFO",SONIQ_DASHBOARD_ENABLED="true",SONIQ_SCHEDULER_SUPPRESS_WARNING="1"

[program:soniq_dashboard]
command=/opt/soniq/venv/bin/soniq dashboard --host=0.0.0.0 --port=8000
directory=/opt/soniq
user=soniq
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/soniq/dashboard.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=SONIQ_DATABASE_URL="postgresql://soniq:password@localhost/soniq_prod",SONIQ_LOG_LEVEL="INFO"

[program:soniq_scheduler]
command=/opt/soniq/venv/bin/soniq scheduler
directory=/opt/soniq
user=soniq
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/soniq/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=SONIQ_DATABASE_URL="postgresql://soniq:password@localhost/soniq_prod",SONIQ_LOG_LEVEL="INFO",SONIQ_SCHEDULING_ENABLED="true"

# Optional: Additional worker for high-priority jobs only
[program:soniq_priority_worker]
command=/opt/soniq/venv/bin/soniq worker --concurrency=2 --queues=urgent
directory=/opt/soniq
user=soniq
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/soniq/priority_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=SONIQ_DATABASE_URL="postgresql://soniq:password@localhost/soniq_prod",SONIQ_LOG_LEVEL="INFO"

# Health check program
[program:soniq_healthcheck]
command=/opt/soniq/scripts/health_check.sh
directory=/opt/soniq
user=soniq
autostart=false
autorestart=false
